# setup/systemd/sync-monitor-gdm.path
####################
# This systemd path file executes the sync-monitor-gdm.service when called
####################

[Unit]
Description=Watch for monitor configuration changes and sync with GDM if changed

[Path]
PathChanged=/home/ridvikpal/.config/monitors.xml
Unit=sync-monitor-gdm.service

[Install]
WantedBy=multi-user.target
